# varlog_tmpfs/varlog_tmpfs
#!/sbin/runscript
# Copyright (c) 2007-2008 Roy Marples <roy@marples.name>
# Released under the 2-clause BSD license.

description="Rsyncs log files at startup/shutdown in case /var/log is a tmpfs"

EXFILE="/etc/varlog_exclude"

depend()
{
	after localmount
}

start()
{
	einfo "Populate /var/log"

	rsync -ar /var/log_permanent/ /var/log/

	eend 0
}

stop()
{
	einfo "Copy chosen files in /var/log to permanent storage"

	rsync -ar /var/log/ /var/log_permanent/ --exclude-from=${EXFILE} --delete-after --delete-excluded

	eend 0
}

